# templates/upstart.conf
description "ceph-metrics-collector"
author "chris <chris@chris-HP-Z400-Workstation>"

start on runlevel [2345]
stop on runlevel [016]

respawn

console log
script
    echo Fake service; sleeping for an hour...
    sleep 360
end script
